# flux-utils/wifi-monitor.conf
[program:wifi-monitor]
command=/opt/root/Flux-Utils/flux-utils/wifi_monitor_launch.zsh
autorestart=true
user=root
